# cosmic-bg-theme.service
[Unit]
Description=Unofficial service for syncing COSMIC theme with the active wallpaper.
After=graphical-session.target

[Service]
Type=simple
ExecStart=/path/to/your/executable
Restart=always
RestartSec=300

[Install]
WantedBy=graphical-session.target